# irys/irys-cli/irys-auto.service
[Unit]
Description=Irys CLI Автоматизация

[Service]
ExecStart=$AUTO_SCRIPT
Restart=always
