# irys/irys-cli/irys-auto.service
[Unit]
Description=Irys CLI Автоматизация

[Service]
ExecStart=/root/.irys/irys_auto.sh
Restart=always
